# kutlak_net/supervisor-website.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/kutlak_net/kutlak_net/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
